# telegrab.service
[Unit]
Description=Telegrab Telegram Archiver
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=telegrab
Group=telegrab
WorkingDirectory=/root/git/telegrab
Environment="PATH=/root/git/telegrab/venv/bin"
ExecStart=/root/git/telegrab/venv/bin/python telegrab.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegrab

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
